# pycheribenchplot/templates/rc-local.sh.jinja
# Auto-generated by cheri-benchplot

. /etc/rc.subr

info "cheri-benchplot rc.local {{ session }}"

# Ensure that we have a kernel file, this is necessary for hwpmc
if [ ! -f /boot/kernel/kernel ]; then
    mkdir -p /boot/kernel
    touch /boot/kernel/kernel
fi

name="benchplot"

load_rc_config $name
: ${benchplot_enable:="NO"}
: ${benchplot_session:=""}
: ${benchplot_user:="root"}

autorun_script="${benchplot_session}/auto-run.sh"

if [ ! -d "${benchplot_session}" ]; then
	err 1 "${benchplot_session} is not a directory."
fi

if [ ! -x "${autorun_script}" ]; then
	err 1 "Autorun script '${autorun_script}' not found or not executable."
fi

info "Autorun benchplot session at ${benchplot_session}"

# -n prevents auto-reboot
(cd ${benchplot_session} && ${autorun_script} -n)

poweroff
